# crontab/statistics_mall.sh.production
#!/usr/bin/env bash
curl http://saas.zhubaoe.cn/index.php/timing/statistics/mall
exit 0